# systemd/hyprswitch.socket
[Unit]
Description=Socket for hyprswitch Daemon

[Socket]
ListenStream=%t/hyprswitch.sock
Accept=false

[Install]
WantedBy=sockets.target